# deploy/systemd/vibe-beat.service
[Unit]
Description=Vibe Sorcery Celery beat (scheduled tasks)
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vibe-sorcery/Vibe-Sorcery-main/backend
EnvironmentFile=/opt/vibe-sorcery/Vibe-Sorcery-main/.env
Environment=PYTHONPATH=/opt/vibe-sorcery/Vibe-Sorcery-main:/opt/vibe-sorcery/Vibe-Sorcery-main/backend
ExecStart=/opt/vibe-sorcery/Vibe-Sorcery-main/venv/bin/celery -A app.celery_app.celery_app beat --loglevel=info
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
